# config/supervisor/supervisord.conf
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0770

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[supervisord]
user=root
nodaemon=true
logfile=/app/logs/supervisord/supervisord.log
pidfile=/var/run/supervisord.pid

[include]
files = /etc/supervisor/conf.d/*.conf

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true

[program:gunicorn]
command=gunicorn -c ./config/gunicorn/gunicorn.conf.py
directory=/app
autorestart=true
